# temp/etc/init/uwsgi.conf
description     "uWSGI Emperor"

start on [2345]
stop on runlevel [!2345]

respawn

env BINPATH=/usr/bin/uwsgi
exec $BINPATH --ini /etc/uwsgi/uwsgi.ini 